# backend/templates/rathole.service.j2
# Auto-generated by Rathole Manager
[Unit]
Description=Rathole Service (role: {{ role }}, config: {{ config_filename }})
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/rathole /etc/rathole/{{ config_filename }}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target